# deploy/systemd/ku-secretary-weather-sync.service
[Unit]
Description=KU Secretary weather + air-quality snapshot sync
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/KUSecretary
EnvironmentFile=%h/KUSecretary/.env
ExecStart=/path/to/KU_Secretary/.venv/bin/kus sync-weather --timeout-seconds 600
TimeoutStartSec=900

[Install]
WantedBy=default.target
