# deploy/systemd/ku-secretary-uclass-poller.service
[Unit]
Description=KU Secretary UClass/Canvas poller
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/KUSecretary
EnvironmentFile=%h/KUSecretary/.env
ExecStart=/path/to/KU_Secretary/.venv/bin/kus sync --all --wait --timeout-seconds 600
TimeoutStartSec=900

[Install]
WantedBy=default.target
