# deploy/systemd/ku-secretary-uclass-poller.timer
[Unit]
Description=Run KU Secretary UClass/Canvas poller every 60 minutes
Requires=ku-secretary-uclass-poller.service

[Timer]
# Conservative interval to avoid abuse detection on KUPID/Canvas.
OnBootSec=5min
OnUnitActiveSec=60min
Unit=ku-secretary-uclass-poller.service

[Install]
WantedBy=timers.target
