# deploy/systemd/ku-secretary-weather-sync.timer
[Unit]
Description=Run KU Secretary weather sync every 30 minutes (offset for jitter)
Requires=ku-secretary-weather-sync.service

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min
RandomizedDelaySec=2min
Unit=ku-secretary-weather-sync.service

[Install]
WantedBy=timers.target
