# deploy/systemd/ku-secretary-briefings.service
[Unit]
Description=KU Secretary morning/evening Telegram briefing dispatch
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/KUSecretary
EnvironmentFile=%h/KUSecretary/.env
ExecStart=/path/to/KU_Secretary/.venv/bin/kus send-briefings
TimeoutStartSec=300

[Install]
WantedBy=default.target
