# deploy/systemd/ku-secretary-briefings.timer
[Unit]
Description=Run KU Secretary briefing dispatch every 15 minutes
Requires=ku-secretary-briefings.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Unit=ku-secretary-briefings.service

[Install]
WantedBy=timers.target
